# crontab
5 * * * *   cd /app/ && python3 send_kindle_new.py > /proc/1/fd/1 2> /proc/1/fd/2
